# deploy/launchd/com.whenever.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  매일 04:30 DB 백업. 설치:
    cp deploy/launchd/com.whenever.backup.plist ~/Library/LaunchAgents/
    # 아래 <REPO_PATH>, <YOU> 를 실제 경로로 치환
    launchctl load ~/Library/LaunchAgents/com.whenever.backup.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.whenever.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string><REPO_PATH>/deploy/scripts/backup-db.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- launchd 는 PATH 가 최소라 docker CLI 경로를 명시 (Apple Silicon: /opt/homebrew/bin) -->
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <key>WorkingDirectory</key>
  <string><REPO_PATH></string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>4</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/<YOU>/whenever-backups/backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/<YOU>/whenever-backups/backup.log</string>
</dict>
</plist>
